# systemd/bot-lnk.service
# Systemd service for bot-lnk with watchdog support
#
# Install:
#   1. Edit paths below (User, WorkingDirectory, ExecStart, EnvironmentFile)
#   2. sudo cp systemd/bot-lnk.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable bot-lnk.service
#   5. sudo systemctl start bot-lnk.service
#
# Monitor:
#   systemctl status bot-lnk --no-pager
#   journalctl -u bot-lnk -f
#
# See docs/SYSTEMD_WATCHDOG.md for details on watchdog support

[Unit]
Description=Lords and Knights Bot (Playwright)
After=network.target solver-lnk.service
Wants=solver-lnk.service

[Service]
# Type=notify enables systemd watchdog support
# The bot will send "I'm alive" pings to systemd
Type=notify
NotifyAccess=main

# TODO: Change to your user/group
User=YOUR_USER
Group=YOUR_USER

# TODO: Change to your repo path
WorkingDirectory=/path/to/bot-lnk

# TODO: Change to your node binary path (run: which node)
ExecStart=/path/to/node dist/index.js

# Watchdog: restart if no ping for 5 minutes
WatchdogSec=300
Restart=on-watchdog
RestartSec=30

# TODO: Uncomment and set path to load .env file
# EnvironmentFile=/path/to/bot-lnk/.env
Environment=NODE_ENV=production

# Logging: stdout to journal, stderr discarded
StandardOutput=journal
StandardError=null

# Basic hardening (Playwright-compatible)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
